# crontab
@reboot sleep 10 && /bin/bash /home/fox282/dashcam.sh > /home/fox282/dashcam.log 2>&1